# posh.service
[Unit]
Description=Printing Over Simple HTTP service

[Service]
ExecStart=/opt/posh

[Install]
WantedBy=multi-user.target
